# systemd/oc-api.service
[Unit]
Description=OpenClaw Flask API (Gunicorn)
After=postgresql@17-main.service redis-server.service
Wants=postgresql@17-main.service redis-server.service

[Service]
Type=exec
User=YOUR_USER
Group=YOUR_USER
WorkingDirectory=/home/YOUR_USER/openclaw-services/backend
EnvironmentFile=/home/YOUR_USER/openclaw-services/.env
ExecStart=/home/YOUR_USER/openclaw-services/backend/venv/bin/gunicorn --bind 127.0.0.1:4001 --workers 1 --threads 4 wsgi:app
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
