# systemd/openclaw-memory.service
[Unit]
Description=OpenClaw Memory Service (FastAPI + bge-m3 + pgvector)
After=postgresql@17-main.service redis-server.service
Wants=postgresql@17-main.service redis-server.service

[Service]
Type=exec
User=YOUR_USER
Group=YOUR_USER
WorkingDirectory=/home/YOUR_USER/agentic-memory/memory-service
ExecStart=/home/YOUR_USER/agentic-memory/memory-service/venv/bin/uvicorn main:app --host 127.0.0.1 --port 18800
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
